# systemd/user/football-commentator.service.d/model.conf.example
[Service]
Environment=OLLAMA_MODEL=gemma4:12b
